# dsa-bot-logrotate.timer
[Unit]
Description=Rotate DSA bot logs hourly

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

