# dsa-bot-logrotate.service
[Unit]
Description=Rotate DSA bot logs

[Service]
Type=oneshot
ExecStart=/ws/vishwsh2-sjc/dsaTelegram/rotate_logs.sh

